# docs/compass.service
[Unit]
Description=Mqtt Compass Service

[Service]
ExecStart=/home/pi/opt/compass/compass
User=pi

[Install]
WantedBy=multi-user.target

